# packaging/flowpulse.service
[Unit]
Description=FlowPulse Intelligence Server
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=flowpulse
Group=flowpulse
WorkingDirectory=/usr/lib/flowpulse/server
EnvironmentFile=/etc/flowpulse/flowpulse.env
ExecStart=/opt/flowpulse/node/bin/node /usr/lib/flowpulse/server/server.js
Restart=always
RestartSec=3

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

# IMPORTANTE: só /var/lib/flowpulse precisa de escrita.
# /usr/lib e /usr/share são read-only (ProtectSystem=strict cuida disso).
# Se o server.js escrever logs/cache locais, use /var/lib/flowpulse.
ReadWritePaths=/var/lib/flowpulse

ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictSUIDSGID=true
RestrictNamespaces=true
RestrictRealtime=true

# Limitar acesso a rede (o serviço só precisa de IPv4/IPv6)
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

StandardOutput=journal
StandardError=journal
SyslogIdentifier=flowpulse

[Install]
WantedBy=multi-user.target
